# add-to-supervisord.conf



[program:xfreerdp]
user=root
priority=100
command=/FreeRDP/xfreerdp.sh
environment=DISPLAY=":1"






